# deploy/vps-netwatch-collector.service
[Unit]
Description=vps-netwatch LAN mihomo collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vps-netwatch
ExecStart=/opt/vps-netwatch/bin/collector -config /opt/vps-netwatch/config.yaml -id lan-proxy-vm
Restart=on-failure
RestartSec=5s
User=root

[Install]
WantedBy=multi-user.target

